# clerk/jetson/clerk-boot-notify.service
[Unit]
Description=Notify that the Jetson Clerk host has booted
After=clerk.service network-online.target
Wants=clerk.service

[Service]
Type=oneshot
EnvironmentFile=/mnt/ssd/clerk-data/clerk.env
ExecStart=/mnt/ssd/open-esquire-verifier/clerk/jetson/clerk-notify.sh "Jetson Clerk host booted; the Clerk is starting." default arrow_up

[Install]
WantedBy=multi-user.target
